# cookbooks/redis/templates/default/redis.upstart.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

description "redis"

start on runlevel [2345]
stop on [!2345]

respawn
respawn limit 10 5

pre-start script
  if [ ! -e <%= node[:redis][:datadir] %> ]; then
    mkdir -p <%= node[:redis][:datadir] %>
  fi
end script

exec sudo -u redis /usr/local/bin/redis-server <%= node[:redis][:config] %>

